# templates/example/fs/etc/systemd/system/firstrun.service
[Unit]
Description=First run

[Service]
Type=oneshot
ExecStart=/etc/firstrun/firstrun.sh

[Install]
WantedBy=multi-user.target

